# deploy/tencent-cvm/momenta-h5.service
[Unit]
Description=Momenta H5 Web
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/momenta/elderly-video-app
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=BACKEND_BASE_URL=http://127.0.0.1:10099
ExecStart=/usr/bin/node /opt/momenta/elderly-video-app/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
